# ubuntu-18.04/etc/rc.local
#!/bin/sh -e
#
# rc.local
#

ip route change default via 172.31.255.254
exit 0
